# deploy/vpntest.service
[Unit]
Description=VPN Proxy Checker
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/proxy-checker
ExecStart=/opt/proxy-checker/venv/bin/python server.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
